# systemd/bark-display.service
[Unit]
Description=Bark Counter - Display
After=bark-detector.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/bark-counter
ExecStart=/home/pi/bark-counter/venv/bin/python -m src.display
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=SDL_VIDEODRIVER=fbcon
Environment=SDL_FBDEV=/dev/fb1

[Install]
WantedBy=multi-user.target
